# root/etc/services.d/fluxbox/run
#!/usr/bin/with-contenv bash

exec fluxbox
