# root/etc/services.d/x11vnc/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc x11vnc \
        -forever \
        -shared
