# root/etc/services.d/xvfb/run
#!/usr/bin/with-contenv bash

exec Xvfb :0 \
    -screen 0 \
    1280x768x24
